# templates/demo-debian-10/files/roles/protected/templates/network/configure_network.service.j2
[Unit]
Description=Add default route for internet access
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/root/networkconf.sh
Type=oneshot
StandardOutput=journal
RemainAfterExit=false

[Install]
WantedBy=multi-user.target
